# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/opt/nginx_config && /usr/sbin/nginx
autorestart=true
